# deploy/synk-watchdog.service
[Unit]
Description=Synk Watchdog (single cycle)
After=network-online.target

[Service]
Type=oneshot
User=root
WorkingDirectory=/root/synk
EnvironmentFile=/root/synk/.env
ExecStart=/root/synk/.venv/bin/python alerts/watchdog.py
StandardOutput=append:/root/synk/logs/watchdog.log
StandardError=append:/root/synk/logs/watchdog.log
